# machines/raspberrypi4/crontab.cron
@reboot sh "$HOME/git/homelab/machines/.current/startup.sh" >~/.log/startup-cron.log 2>&1
