# docker-images/data-optimization/crontab.cron
00 */4 * * * sh /homelab/cron-wrapper.sh >>/homelab/logs/cron.log 2>&1
