# machines/odroid-h3/crontab.cron
@reboot sh "$HOME/git/homelab/.current-machine/startup.sh"
