# machines/odroid-h3/crontab.cron
@reboot sh "$HOME/git/homelab/machines/.current/startup.sh"
