# servers/raspberry-pi-4b-2g/crontab.cron
@reboot sh "$HOME/git/homelab/servers/.current/startup.sh" >"$HOME/.log/startup-cron.log" 2>&1
