# docker-images/external/renovatebot/crontab.cron
30 1,3,5 * * * sh /homelab/cron-wrapper.sh >>/homelab/logs/cron.log 2>&1
