# docker-images/external/renovatebot/crontab.cron
0 1 * * * sh /homelab/cron-wrapper.sh >>/homelab/logs/cron.log 2>&1
