# servers/raspberry-pi-4b-4g/crontab.cron
@reboot sh "$HOME/startup.sh" >"$HOME/.log/startup-cron.log" 2>&1
