# machines/odroid-h3/unifi-network-application/backup/schedule.cron
* * * * * sh /app/run.sh >>/log/cron.log 2>&1
