# docker-images/custom/docker-build/crontab.cron
0 2 * * * sh /homelab/cron-wrapper.sh >>/homelab/logs/cron.log 2>&1
