# docker-images/custom/web-automation/crontab-other.cron
00 04 * * * sh /homelab/cron-wrapper.sh >>/homelab/logs/cron.log 2>&1
