# docker-images/custom/web-automation/crontab-other.cron
00 04 * * * sh /app/cron-wrapper.sh >>/log/cron.log 2>&1
