# machines/raspberry-pi-zero-2w-2/crontab.cron
@reboot sh "$HOME/git/homelab/machines/.current/startup.sh" >"$HOME/.log/startup-cron.log" 2>&1
