# components/webui-backup/crontab.cron
30 03 * * * sh /app/run.sh >>/log/cron.log 2>&1
