# servers/raspberry-pi-3b/config/crontab.cron
@reboot sh "$HOME/config/startup.sh" >"$HOME/.log/startup-cron.log" 2>&1
30 5 * * * sh "$HOME/config/update.sh" >"$HOME/.log/update-cron.log" 2>&1
